# supervisord/supervisord.ini
[unix_http_server]
file=/root/PycharmProjects/stock/supervisord/supervisor.sock ; (the path to the socket file)
chmod=0700                  ; socket file mode (default 0700)
chown=umail:umail           ; socket file uid:gid owner

[supervisord]
logfile=/root/PycharmProjects/stock/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=5MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=2           ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/root/PycharmProjects/stock/supervisord/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
user=root                  ; (default is current user, required if root)
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///root/PycharmProjects/stock/supervisord/supervisor.sock ; use a unix:// URL  for a unix socket
prompt=u-mail               ; cmd line prompt (default "supervisor")

[group:data_update]
programs=daily, weekly, monthly

[program:daily]
command                 = /root/anaconda2/bin/python /root/PycharmProjects/stock/task_daily_update.py
directory               = /root/PycharmProjects/stock
autostart               = true
autorestart             = true
redirect_stderr         = true
stdout_logfile          = /root/PycharmProjects/stock/log/task_daily_update.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 2

[program:weekly]
command                 = /root/anaconda2/bin/python /root/PycharmProjects/stock/task_weekly_update.py
directory               = /root/PycharmProjects/stock
autostart               = true
autorestart             = true
redirect_stderr         = true
stdout_logfile          = /root/PycharmProjects/stock/log/task_weekly_update.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 2

[program:monthly]
command                 = /root/anaconda2/bin/python /root/PycharmProjects/stock/task_monthly_update.py
directory               = /root/PycharmProjects/stock
autostart               = true
autorestart             = true
redirect_stderr         = true
stdout_logfile          = /root/PycharmProjects/stock/log/task_monthly_update.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 2

